# cookbooks/default/files/default/systemd/chef-client.timer
[Unit]
Description=Chef Client Timer
[Timer]
OnUnitActiveSec=30m
Unit=chef-client.service
[Install]
WantedBy=multi-user.target
